# ai-agent.service
[Unit]
Description=AI Agent — Claude-powered autonomous agent service
Documentation=https://github.com/abusultancom/-ai-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/user/-ai-agent
EnvironmentFile=/home/user/-ai-agent/.env
ExecStart=/usr/local/bin/python3 /home/user/-ai-agent/orchestrator.py serve
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-agent

# Resource limits — adjust to taste
MemoryMax=4G
TasksMax=256

[Install]
WantedBy=multi-user.target
